# install/etc/profile.d/inist-tools.sh
export PATH="$PATH:/opt/inist-tools"
. /opt/inist-tools/inistrc
