# install/etc/profile.d/inist-tools.sh
#!/usr/bin/env bash
export PATH=$PATH:/opt/istex-tools
